# util/build-deb/prerm
#!/bin/sh

if [ "$1" = remove ]; then
    systemctl disable --now unlockr.service || true
fi

exit 0
